# contrib/systemd/lockcd.service.in
[Unit]
Description=lockc daemon
After=network-online.target

[Service]
Type=oneshot
ExecStart={{ bindir }}/lockcd
StandardOutput=journal

[Install]
WantedBy=multi-user.target
