# contrib/systemd/lockcd.service.in
[Unit]
Description=lockc daemon
After=network-online.target

[Service]
Type=simple
Restart=always
RestartSec=1
ExecStart={{ bindir }}/lockcd
StandardOutput=journal

[Install]
WantedBy=multi-user.target
